# etc/profile.d/cl_path.sh
# This script adds cumulus binary paths to PATH

PATH=$PATH:/usr/cumulus/bin
export PATH

